# s6-overlay/s6-rc.d/haproxy/run
#!/bin/bash
exec 2>&1
exec haproxy -f /etc/haproxy/haproxy.cfg -db
