# s6-overlay/s6-rc.d/tor/finish
#!/bin/bash
killall tor 2>/dev/null || true
