# s6-overlay/s6-rc.d/privoxy/run
#!/bin/bash
exec 2>&1
exec privoxy --no-daemon /etc/privoxy/config
